# bot.conf
[program:jadvaluzbot]
directory=/home/ubuntu/JadvalUzbot
command=python3 app.py

autostart=true
autorestart=true
enviroment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/JadvalUzbot/logfile_err.log
stdout_logfile=/home/ubuntu/JadvalUzbot/logfile.log